# ansible/roles/sync/templates/systemd/dotfiles-sync.service.j2
[Unit]
Description=Dotfiles GitOps Sync
Documentation=file://{{ dotfiles_dir }}/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles-sync
StandardOutput=journal
StandardError=journal
TimeoutSec=120
SuccessExitStatus=0

[Install]
WantedBy=default.target
