# ansible/roles/nvim/templates/nvim-config-sync.timer.j2
[Unit]
Description=Neovim Config GitOps Sync Timer
Requires=nvim-config-sync.service

[Timer]
# Fire shortly after login so the first nvim session of a new day picks up
# any upstream changes made on another machine
OnBootSec=3min
# Daily sync — nvim config changes are less frequent than shell config
OnCalendar=daily
# Catch up if the machine was off during the scheduled window
Persistent=true
# Spread starts within a 5-minute window to avoid hammering origin at midnight
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
