# ansible/roles/ai-tools/templates/launchd/com.ai-config.sync.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-config.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ ai_sync_script_dest }}</string>
    </array>

    <!-- 30-minute interval in seconds to match the systemd timer cadence -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Fire once on load so a fresh install picks up upstream commits -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{ ai_state_dir }}/logs/sync.log</string>

    <key>StandardErrorPath</key>
    <string>{{ ai_state_dir }}/logs/sync-error.log</string>

    <!-- Throttle restarts on repeated failure (SSH broken, no network, etc.) -->
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
